# supervisord-rsyslog.conf
[program:rsyslog]
command=/configure-and-run-rsyslog.sh
numprocs=1
autostart=true
autorestart=false
